# src/crontab
* * * * * /usr/local/bin/python /app/script.py >>/tmp/out.log 2>&1
